# src/static/static/home/yi-hack-v5/etc/crontabs/root
1 1 * * * /tmp/sd/yi-hack-v5/script/check_update.sh
